# systemd/backup.timer
[Unit]
Description=Run weekly backup (Sunday 03:30)

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
RandomizedDelaySec=5m
Unit=backup.service

[Install]
WantedBy=timers.target
