# systemd/maintenance.timer
# /etc/systemd/system/maintenance.timer
[Unit]
Description=Run weekly maintenance (Saturdays 04:00)

[Timer]
OnCalendar=Sat *-*-* 04:00:00
Persistent=true
RandomizedDelaySec=5m
Unit=maintenance.service

[Install]
WantedBy=timers.target
