# systemd/oneclick-status.service
[Unit]
Description=OneClick daily health check

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '/usr/local/bin/oneclick-status | tee -a /var/log/oneclick-status.log'
