# systemd/wg-autoswitch.service
[Unit]
Description=WireGuard auto-switch based on streaming metrics
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/osmc-oneclick/phases/31_vpn_autoswitch.sh
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
