# systemd/wg-autoswitch.service
[Unit]
Description=WireGuard auto-switch based on streaming metrics
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/osmc-oneclick/phases/31_vpn_autoswitch.sh
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=5

[Install]
WantedBy=multi-user.target
