# firstboot/oneclick-firstboot.service
[Unit]
Description=OneClick first boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /boot/firstboot.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
